# build_files/thunderbolt-extension/silverletter-thunderbolt-disconnect.service
[Unit]
Description=Power down Thunderbolt after the external device disconnects

[Service]
Type=oneshot
ExecStart=/usr/libexec/silverletter-thunderbolt-control hotplug-release
TimeoutStartSec=45
