# build_files/thunderbolt-extension/silverletter-thunderbolt-disconnect.service
[Unit]
Description=Power down Thunderbolt after the external device disconnects
StartLimitIntervalSec=600
StartLimitBurst=2

[Service]
Type=oneshot
ExecStart=/usr/libexec/silverletter-thunderbolt-control hotplug-release
# Falcon Ridge sysfs removal can block for tens of seconds per bridge.
TimeoutStartSec=180
Restart=on-failure
RestartSec=5
