# build_files/power/aspm-tune-resume.service
[Unit]
Description=ASPM tuning after resume from suspend
After=suspend.target hibernate.target hybrid-sleep.target suspend-then-hibernate.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/aspm-tune.sh
RemainAfterExit=no

[Install]
WantedBy=suspend.target hibernate.target hybrid-sleep.target suspend-then-hibernate.target
