# build_files/mactahoe/mactahoe-firefox-setup.service
[Unit]
Description=MacTahoe Desktop and Firefox Theme Setup
PartOf=graphical-session.target

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/libexec/mactahoe-firefox-setup.sh

[Install]
WantedBy=graphical-session.target
