# build_files/toshy/toshy-first-login-setup.service
[Unit]
Description=Toshy First-Login Setup
ConditionPathExists=!%h/.config/toshy/.image-setup-done
ConditionPathExists=!%h/.config/toshy/toshy_config.py
After=graphical-session.target
Requires=graphical-session.target

[Service]
Type=simple
RemainAfterExit=yes
# Run the terminal directly here — stays alive until the window closes
ExecStart=/usr/libexec/toshy-first-login-launch.sh
TimeoutStartSec=600
Restart=no

[Install]
WantedBy=graphical-session.target